# agent-runtime/scripts/podman/builderforce.container.in
# OpenAgentNode gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into openagentNode's ~/.config/containers/systemd/
# {{OPENAGENT_NODE_HOME}} is replaced at install time.

[Unit]
Description=OpenAgentNode gateway (rootless Podman)

[Container]
Image=openagentNode:local
ContainerName=openagentNode
UserNS=keep-id
Volume={{OPENAGENT_NODE_HOME}}/.openagentNode:/home/node/.openagentNode
EnvironmentFile={{OPENAGENT_NODE_HOME}}/.openagentNode/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
